# deploy/steep-digest.service
[Unit]
Description=Steep newsletter digest (Gmail -> LLM -> Telegram + Gmail)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=steep
Group=steep
WorkingDirectory=/opt/steep
EnvironmentFile=-/etc/steep/digest.env
ExecStart=/opt/steep/.venv/bin/steep-digest run

# Hardening (relax if you hit permission issues)
NoNewPrivileges=true
PrivateTmp=true
